# ai-pdf-screen.service
[Unit]
Description=AI PDF Brain — fullscreen insight display
# Wait for the graphical session (X on the MHS35 TFT) to be up.
After=graphical.target ollama.service
Wants=ollama.service
Requires=graphical.target

[Service]
Type=simple
User=pi
Group=pi
Environment=HOME=/home/pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/ai-pdf
ExecStart=/home/pi/ai-pdf/start.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SupplementaryGroups=video render input

[Install]
WantedBy=graphical.target
